# tools/paperclip_watchdog/com.tubeslave.automixer-paperclip-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tubeslave.automixer-paperclip-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>python3</string>
    <string>/Users/dmitrijvolkov/AUTO-MIXER-Tubeslave-main/tools/paperclip_watchdog/automixer_paperclip_watchdog.py</string>
    <string>loop</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/dmitrijvolkov/AUTO-MIXER-Tubeslave-main</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/dmitrijvolkov/AUTO-MIXER-Tubeslave-main/logs/automixer_paperclip_watchdog.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/dmitrijvolkov/AUTO-MIXER-Tubeslave-main/logs/automixer_paperclip_watchdog.launchd.err.log</string>
</dict>
</plist>
